# root/etc/s6-overlay/s6-rc.d/nzbhydra2/run
#!/command/with-contenv sh

cd /opt/nzbhydra2
exec s6-setuidgid appuser python3 /opt/nzbhydra2/nzbhydra2wrapperPy3.py
