# root/etc/s6-overlay/s6-rc.d/nzbhydra2-fixattrs/up
/etc/s6-overlay/s6-rc.d/nzbhydra2-fixattrs/run
